# systemd/station-announcer.timer
[Unit]
Description=Run station-announcer.service every 10 minutes

[Timer]
OnCalendar=*:0/10
