# systemd/station-announcer.timer
[Unit]
Description=Run station-announcer.service every 25 minutes

[Timer]
OnCalendar=*:0/25
