# systemd/station-announcer-timed.timer
[Unit]
Description=Run station-announcer-timed.service every 27 minutes

[Timer]
OnCalendar=*:0/27
